# laravel-app/config/supervisor/laravel-reverb.conf
[program:laravel-reverb]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/html/laravel-app/artisan reverb:start --host=0.0.0.0 --port=8080
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/html/laravel-app/storage/logs/reverb.log
stopwaitsecs=10
